# octoprint_slim/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: OctoPrint (Slim)
# Starts server to serve the app
# ==============================================================================
bashio::log.info "Run OctoPrint (Slim).."

cd /opt/panel \
    || bashio::exit.nok "Could not change into app directory"

copy_data() {
    if [ ! -d /data/python ]; then
        bashio::log.info "Copy authentication information"
        cp -R /root/python /data/
        cp -R /config/octoprint/users.yaml /data/users.yaml
    fi
}

copy_data

octoprint serve --iknowwhatimdoing --basedir /data --config /data/options.json --host 0.0.0.0 --port 5000
