# octoprint_slim/rootfs/etc/services.d/octoprint/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: OctoPrint (Slim)
# Starts server to serve the app
# ==============================================================================
bashio::log.info "Run OctoPrint (Slim).."

copy_data() {
    if [ ! -d /data/python ]; then
        bashio::log.info "Copy authentication information"
        cp -R /root/python /data/
        cp -R /config/octoprint/users.yaml /data/users.yaml
    fi
}

copy_data

# Extract OctoPrint server configuration
jq ".options.octoprint" /data/options.json > /data/octoprint.json

# Start OctoPrint server
octoprint serve --iknowwhatimdoing --basedir /data --config /data/octoprint.json --host 0.0.0.0 --port 5000
